# systemd/borg-backup.timer
[Unit]
Description=Automated Borg Backup (Timer)

[Timer]
OnCalendar=02:30

[Install]
WantedBy=timers.target