# systemd/borg-backup.service
[Unit]
Description=Automated Borg Backup (Service)

[Service]
Type=oneshot
ExecStart=/home/user/borgbackup/backup.sh